# playbooks/service/templates/todo-app.service.j2
[Unit]
Description=Todo App Server

[Service]
EnvironmentFile=/etc/environment
ExecStart=/usr/bin/todo-app

[Install]
WantedBy=multi-user.target